# las-nubes.service.template
[Unit]
Description=Las Nubes - Automation Control
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%INSTALL_DIR%
Environment=NODE_ENV=production
EnvironmentFile=%INSTALL_DIR%/.env
ExecStart=/usr/bin/node %INSTALL_DIR%/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=las-nubes

[Install]
WantedBy=multi-user.target
